# rootfs/etc/profile.d/nobody_no.sh
if (( $EUID == 65534 )); then
    sudo su -l delta
    exit
fi